# rootfs/etc/s6-rc/image/service/getty/run
#!/bin/execlineb -P
/sbin/getty -l /usr/bin/autologin 38400 tty1
