# rootfs/etc/profile.d/banner.sh
#!/bin/sh

echo "Redpill Linux version $(cat /etc/version)"
docker -v
